# postfix/service/postfix/run
#!/bin/bash
set -e

# Avoid  warning: smtputf8_enable is true, but EAI support is not compiled in
echo "smtputf8_enable = no" >> /etc/postfix/main.cf

# Do we want to modify the config first with the script?
[ -f /etc/service/postfix/run.config ] && source /etc/service/postfix/run.config

if [[ ! -z "$MAILNAME" ]]; then
	echo "$MAILNAME" > /etc/mailname
	postconf -e myhostname="$MAILNAME"
fi

if [[ ! -z "$MY_NETWORKS" ]]; then
	postconf -e mynetworks="$MY_NETWORKS"
fi

if [[ ! -z "$MY_DESTINATION" ]]; then
	postconf -e mydestination="$MY_DESTINATION"
fi

if [[ ! -z "$ROOT_ALIAS" ]]; then
	if [[ -f /etc/aliases ]]; then
		sed -i '/^root:/d' /etc/aliases
	fi
	echo "root: $ROOT_ALIAS" >> /etc/aliases
	newaliases
fi

if [[ ! -z "$RELAY" ]]; then
	# setup the relay
	echo "relay_host = $RELAY" >> /etc/postfix/main.cf
fi

if [[ ! -z "$SASL_AUTH" ]]; then
	# setup tls
	echo -e "smtp_sasl_auth_enable = yes\nsmtp_sasl_password_maps = hash:/etc/postfix/sasl_passwd\nsmtp_sasl_security_options = noanonymous\nsmtp_use_tls = yes\ninet_interfaces = 127.0.0.1" >> /etc/postfix/main.cf

	# generate the SASL password map
	echo "$RELAY $SASL_AUTH" > /etc/postfix/sasl_passwd

	# generate a .db file
	postmap /etc/postfix/sasl_passwd

	# cleanup
	rm /etc/postfix/sasl_passwd

	# set permissions
	chmod 600 /etc/postfix/sasl_passwd.db
fi

exec /usr/lib/postfix/master -c /etc/postfix -d 2>&1
tail -F /var/log/mail.log
